# server/reservoir.service
[Unit]
Description=Reservoir Service

[Service]
ExecStart=PATH_TO_EXEC

[Install]
WantedBy=multi-user.target